# gigabox/evremap/laptop-remap.service
[Service]
WorkingDirectory=/
ExecStart=bash -c "/usr/bin/evremap remap /etc/laptop.toml -d 0"
Restart=always

[Install]
WantedBy=multi-user.target
